# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/webui.conf
[program:webui]
command=/opt/ai-dock/bin/supervisor-webui.sh
process_name=%(program_name)s
numprocs=1
directory=/root
priority=1500
autostart=true
startsecs=0
startretries=3
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/webui.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/webui.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=1
environment=PROC_NAME="%(program_name)s"